# gf/protected/extensions/crontab/crontabs/my_crontab
* 2 * 6 5 php C:/wamp/www/yiiPrueba/GfoxV3/yiicmd.php test 'datetime
